# deploy/bsg-beta.service
# systemd-Unit für die BETA-Stufe (Hetzner).
# Installation:  sudo cp deploy/bsg-beta.service /etc/systemd/system/
#                sudo systemctl daemon-reload && sudo systemctl enable --now bsg-beta
#
# BSG_DEV=1: Beta erlaubt devCode-Login + /api/test/reset (für QA/Contract-Smoke).
# Deshalb MUSS die Beta-Subdomain in Caddy per basicauth geschützt sein (s. Caddyfile).
[Unit]
Description=BSG Benninghausen – Beta (Node)
After=network.target

[Service]
Type=simple
User=bsgdeploy
Group=bsgdeploy
WorkingDirectory=/var/www/bsg-beta
ExecStart=/usr/bin/node /var/www/bsg-beta/server/index.mjs
Environment=PORT=8081
Environment=BSG_DEV=1
Environment=BSG_STATIC=1
Environment=BSG_SECURE_COOKIES=1
Restart=always
RestartSec=2
# Härtung
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/www/bsg-beta
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
